# systemd/github-ntfy.service
[Unit]
Description=Forward GitHub notifications to ntfy
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=github-ntfy
Group=github-ntfy
EnvironmentFile=/etc/github-ntfy.env
Environment=HOME=/var/lib/github-ntfy
ExecStart=/usr/local/bin/github-ntfy poll
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/github-ntfy
UMask=0077
