# systemd/user/github-ntfy.timer
[Unit]
Description=Poll GitHub notifications for ntfy forwarding (user)

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min
AccuracySec=1s
RandomizedDelaySec=0
Unit=github-ntfy.service

[Install]
WantedBy=timers.target
