# systemd/user/github-ntfy.service
[Unit]
Description=Forward GitHub notifications to ntfy (user)
After=network-online.target

[Service]
Type=oneshot
EnvironmentFile=%h/.config/github-ntfy.env
Environment=STATE_DIR=%h/.local/state/github-ntfy
ExecStart=%h/.local/bin/github-ntfy poll
NoNewPrivileges=yes
PrivateTmp=yes
UMask=0077
